# certBOOT.service
[Unit]
Description=Certbot starter

[Service]
ExecStart=/usr/bin/certBOOT.sh

[Install]
WantedBy=multi-user.target